# etc/systemd/system/msmtp.timer
[Unit]
Description=Run msmtp queue

[Timer]
OnBootSec=10min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
